# mac/contact-it-recurring-prompt/Library/LaunchAgents/com.github-angela-d.contactIT.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>the_title</key>
		<string>Contact IT</string>
		<key>the_subtitle</key>
		<string></string>
		<key>bundle_id</key>
		<string>com.github-angela-d.contactIT</string>
	</dict>
	<key>Label</key>
	<string>com.github-angela-d.contactIT</string>
	<key>Program</key>
	<string>/usr/local/sbin/contact-it-recurring-prompt.sh</string>
	<key>StartInterval</key>
	<integer>1800</integer>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
